# scripts/cron-levels.cron
# Rumb Comarcal cron (00:01)
1 0 * * * /Users/pau/Desktop/GitHub/Rumb/Rumb/scripts/cron-levels.sh daily
